# modules/stop/guest/stop-done/mkosi.extra/usr/local/lib/systemd/system/stop-done.service
[Unit]
Description=Stop Barrier
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/true
RemainAfterExit=yes

[Install]
WantedBy=stop.target
