# modules/build/host/targets/mkosi.extra/usr/local/lib/systemd/system/launch.target
# Target reached when the snp guest has launched.
# Job services should use: WantedBy=launch.target

[Unit]
Description=Launch Stage

After=system.target
Requires=system.target

Wants=launch-done.service
After=launch-done.service
