# modules/launch/host/guest-measurement/mkosi.extra/usr/local/lib/systemd/system/calculate-measurement.service
[Unit]
Description=Calculate guest measurement
DefaultDependencies=no
After=system.target
Wants=system.target

[Service]
Type=oneshot
ExecStart=/usr/local/lib/scripts/guest_measurement.sh
StandardOutput=journal+console
StandardError=journal+console
LogExtraFields="SEV_VERSION=3.0.0-0" "SNPHOST_TEST=3.0.0-0"
