# modules/test/guest/attestation-result/mkosi.extra/usr/local/lib/systemd/system/attestation-result.service
[Unit]
Description=Fetch SNP Guest Attestation Status inside the guest
DefaultDependencies=no
After=attestation-workflow.service
Wants=attestation-workflow.service

[Service]
Type=oneshot
ExecStart=jq '.' /usr/local/lib/attestation_status
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5s
